# src/photo-svn106032/bin/everynight.cron
30 22 * * 0-4 /usrdevel/s1/heidi/photoinstall.sh
